# deploy/pi3/lte-swd-cert-renew.timer
[Unit]
Description=Daily check for LTE_SWD TLS certificate renewal

[Timer]
OnCalendar=*-*-* 03:25:00
RandomizedDelaySec=20m
Persistent=true

[Install]
WantedBy=timers.target
